# templates/default/upstart.conf.erb
# <%= @component_name %> Service

description     "<%= @component_name %> Server"
author          "Trotter Cashion <cashion@gmail.com>"

start on (net-device-up
          and filesystem
      and runlevel [2345])
stop on runlevel [016]

respawn

env PATH="<%= @path %>:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
export PATH

exec <%= @binary %> -c <%= @config_file %>
